# noninteractive/.crontab
#
# $Id: .crontab 4221 2012-07-11 14:18:05Z KevinAnLipscomb $
#
00 02 * * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/p4o/noninteractive/milestone_sweep.aspx
15 02 * * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/p4o/noninteractive/daily_tasks.aspx
