# trunk/p4o/noninteractive/.crontab
#
# $Id$
#
00 02 * * * wget --output-document=/dev/null https://sslws038.alentus.com/frompaper2web/p4o/noninteractive/milestone_sweep.aspx
15 02 * * * wget --output-document=/dev/null https://sslws038.alentus.com/frompaper2web/p4o/noninteractive/daily_tasks.aspx
